# init/supervisor-config.conf
[program:transcoder]
command=/usr/local/bin/transcoder.py
user=root
autorestart=true
